# assets/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 1MB
directory=/
command=/usr/sbin/httpd -DFOREGROUND
