# assets/var/spool/cron/crontabs/www
*/5 *	* * *	php8 -f /var/www/html/cron.php
