# assets/var/spool/cron/crontabs/www
*/5 *	* * *	php8 /var/www/html/cron.php
